# main.cron
30 3 * * * /home/deaian/trsh/project/main.py
